# supervisord.conf
[supervisord]
nodaemon=true

[program:server.py]
command=/usr/bin/python /my_application/server.py

[program:sshd]
command=/usr/sbin/sshd -D
